# systemd/megasync.service
[Unit]
Description=megasync client

[Service]
ExecStart=/usr/bin/megasync

[Install]
WantedBy=default.target
