# ansible/roles/common/templates/95-monitor-hotplug.rules.template
KERNEL=="card0", SUBSYSTEM=="drm", ENV{LC_ALL}="en_US.utf-8", ENV{LANG}="en_US.utf-8", ENV{DISPLAY}=":0", ENV{XAUTHORITY}="{{ home }}/.Xauthority", RUN+="{{ home }}/scripts/myrandr"
